# init/dwgd.socket
[Unit]
Description=dwgd

[Socket]
ListenStream=/run/docker/plugins/dwgd.sock

[Install]
WantedBy=sockets.target
